# cultura/mak_plataforma/mak-hub.service
; Espejo del unit vivo en el box: ~/.config/systemd/user/mak-hub.service
; (2026-07-18: unit existia pero DISABLED -- el hub solo vivia por watchdog;
; fix: WorkingDirectory pineado + enable --now. Misma clase de fragilidad
; de reboot que tenia mak-codex.)
;
; Instalar en el box:
;   cp mak-hub.service ~/.config/systemd/user/mak-hub.service
;   systemctl --user daemon-reload
;   systemctl --user enable --now mak-hub.service
[Unit]
Description=MAK hub (plataforma :8900)
[Service]
WorkingDirectory=%h/plataforma
ExecStart=/usr/bin/python3 %h/plataforma/hub.py
Restart=always
RestartSec=5
[Install]
WantedBy=default.target
